# config/nacho-cove.cron
# Delay a random time between 0 and 360 minutes (6 hours)
RANDOM_DELAY=360
MAILTO=janv@nachocove.com,azimo@nachocove.com
1 0 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/T3RedShiftLoader.py --config /home/ec2-user/src/nacho/Telemetry/config/dev_t3_redshift.json --period daily --email
1 0 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/T3RedShiftLoader.py --config /home/ec2-user/src/nacho/Telemetry/config/alpha_t3_redshift.json --period daily --email
1 0 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/T3RedShiftLoader.py --config /home/ec2-user/src/nacho/Telemetry/config/beta_t3_redshift.json --period daily --email
0 2 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/T3LogReport.py --config /home/ec2-user/src/nacho/Telemetry/config/alpha_t3_redshift.json --period daily --email
0 2 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/T3LogReport.py --config /home/ec2-user/src/nacho/Telemetry/config/dev_t3_redshift.json --period daily --email
0    6 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/monitor_dev.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
0    6 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/monitor_alpha.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
0    6 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/monitor_prod.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
0    6 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/monitor_beta1.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
*/5 * * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/support_prod.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --before now --period 300 --email
*/5 * * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/support_beta1.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --before now --period 300 --email
*/5 * * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/support_alpha.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --before now --period 300 --email
0    5 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/email_report_prod.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
0    5 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/email_report_beta1.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --period daily --email
0    */6 * * * ec2-user /usr/bin/python2.7 /home/ec2-user/src/nacho/Telemetry/scripts/monitor.py --config /home/ec2-user/src/nacho/Telemetry/config/monitor_prod_pinger.cfg --email-config /home/ec2-user/src/nacho/Telemetry/config/aws_ses_email.cfg --after last --before now-3hrs --period 21600 --email
